# com.hotel-common.guardrails.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hotel-common.guardrails</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>node</string>
        <string>/Users/kaneko/hotel-common/scripts/continuous-guardrails-monitor.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/kaneko/hotel-common</string>
    
    <key>StandardOutPath</key>
    <string>/Users/kaneko/hotel-common/logs/guardrails-background.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/kaneko/hotel-common/logs/guardrails-error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist> 